# includes.container/usr/lib/systemd/user/autostart.service
[Unit]
Description=Hyprland Wayland compositor
Wants=xdg-desktop-autostart.target

[Service]
Type=oneshot
ExecStart=/usr/bin/true
RemainAfterExit=yes